# systemd/roku-watchdog.service
[Unit]
Description=Roku-E8C3 health watchdog
Documentation=https://github.com/roku-e8c3
After=roku-dashboard.service roku-netapply.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/roku-gateway/roku-gateway.env
WorkingDirectory=/opt/roku-gateway/backend
ExecStart=/opt/roku-gateway/venv/bin/python -m app.watchdog
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
