# systemd/roku-recovery.service
[Unit]
Description=Roku-E8C3 boot recovery check
Documentation=https://github.com/roku-e8c3
After=roku-dashboard.service roku-netapply.service
Wants=roku-dashboard.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/roku-gateway/roku-gateway.env
WorkingDirectory=/opt/roku-gateway/backend
# Let the dashboard finish starting before the health check runs.
ExecStartPre=/bin/sleep 8
ExecStart=/opt/roku-gateway/venv/bin/python -m app.recovery
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
