# systemd/roku-firewall.service
[Unit]
Description=Roku-E8C3 rebuild firewall ruleset from database
After=roku-netapply.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/roku-gateway/roku-gateway.env
WorkingDirectory=/opt/roku-gateway/backend
ExecStart=/opt/roku-gateway/venv/bin/python -m app.netapply firewall

[Install]
WantedBy=multi-user.target
