# systemd/sand-display.service
[Unit]
Description=Roku-E8C3 e-paper status display refresh
Documentation=https://github.com/Reen06/Sand-OS
After=sand-netapply.service sand-dashboard.service

[Service]
Type=oneshot
User=gateway
EnvironmentFile=-/etc/sandos/sandos.env
WorkingDirectory=/home/gateway/roku-gateway/backend
Environment=PYTHONPATH=/home/gateway/roku-gateway/backend
ExecStart=/usr/bin/python3 -m app.display_update
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
