# systemd/sand-dashboard.service
[Unit]
Description=Roku-E8C3 dashboard backend
Documentation=https://github.com/roku-e8c3
After=sand-netapply.service

[Service]
Type=simple
User=sand
Group=sand
EnvironmentFile=-/etc/sandos/sandos.env
WorkingDirectory=/opt/sandos/backend
ExecStart=/opt/sandos/venv/bin/python -m app.main
Restart=always
RestartSec=3
# Allow binding port 80 as a non-root user.
AmbientCapabilities=CAP_NET_BIND_SERVICE
# The dashboard must keep sudo access to the privileged helpers.
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/var/lib/sandos /var/log/sandos

[Install]
WantedBy=multi-user.target
