# systemd/sand-recovery.service
[Unit]
Description=Roku-E8C3 boot recovery check
Documentation=https://github.com/roku-e8c3
After=sand-dashboard.service sand-netapply.service
Wants=sand-dashboard.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sandos/sandos.env
WorkingDirectory=/opt/sandos/backend
# Let the dashboard finish starting before the health check runs.
ExecStartPre=/bin/sleep 8
ExecStart=/opt/sandos/venv/bin/python -m app.recovery
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
